# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx

[program:xweb]
command=/app/xweb -c /app/conf.json
