# supervisord.conf
[supervisord]
nodaemon=true

[program:filebrowser]
command=nginx

[program:xweb]
command=/app/xweb -c /app/conf.json
