# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true

[program:xweb]
command=/app/xweb -c /app/conf.json
autostart=true
autorestart=true
