# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g daemon off

[program:xweb]
command=/app/xweb -c /app/conf.json
